# scripts/launchd-wrappers/com.kuro.hn-ai-trend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kuro.hn-ai-trend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/user/Workspace/mini-agent/scripts/launchd-wrappers/hn-ai-trend.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/user/Workspace/mini-agent/memory/logs/hn-trend-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/user/Workspace/mini-agent/memory/logs/hn-trend-launchd.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/user/Workspace/mini-agent</string>
</dict>
</plist>
